# provisioning/vmc/files/vmc.service.j2
[Unit]
Description=VMC
After=network.target

[Service]
User={{ ansible_user_id }}
ExecStart={{ ansible_user_dir }}/vmc/src/vmc
Restart=on-failure

[Install]
WantedBy=multi-user.target
